# base_external/rootfs_overlay/etc/init.d/S98lddmodules
#!/bin/sh

### BEGIN INIT INFO
# Provides:          S98lddmodules
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Load and unload ldd kernel modules
### END INIT INFO

case "$1" in
    start)
        echo "Loading scull driver..."
        insmod /lib/modules/scull.ko

        echo "Loading faulty driver..."
        insmod /lib/modules/faulty.ko

        echo "Loading hello module..."
        modprobe hello
        ;;
    
    stop)
        echo "Removing hello module..."
        rmmod hello

        echo "Removing faulty driver..."
        rmmod faulty

        echo "Removing scull driver..."
        rmmod scull
        ;;
    
    restart)
        $0 stop
        $0 start
        ;;

    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
        ;;
esac

exit 0

